# storyteller.service
[Unit]
Description=Story Maker - Storytelling Workshop Tool
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/storyteller
ExecStart=/home/pi/.nvm/versions/node/v22.22.3/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storyteller

[Install]
WantedBy=multi-user.target
